# deploy/systemd/evinfra-backend.service
[Unit]
Description=EV Infra Backend (Flask via Gunicorn)
After=network.target

[Service]
Type=simple
User=evinfra
Group=evinfra
# Adjust to your deployment path
WorkingDirectory=/opt/evinfra/backend
# Environment file to hold secrets and DATABASE_URL
EnvironmentFile=/etc/evinfra/backend.env
# Use project virtual environment gunicorn
ExecStart=/opt/evinfra/backend/venv/bin/gunicorn --workers 3 --bind 127.0.0.1:5000 "app:create_app()"
Restart=on-failure
RestartSec=5
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
